# automation/launchd/com.optionharvest.history.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.optionharvest.history</string>

  <!-- Launches through ThetaHarvest.app (a signed wrapper) instead of python3 directly, so
       Login Items & Extensions shows "Theta Harvest" + the brand icon rather than
       "Python Software Foundation" (the python.org code-sign identity). The app's executable
       just sets env and execs `python3 -m automation.run_history_update`.
       Rebuild the app with: zsh automation/launchd/build_app.sh (or run install.sh).
       Full Disk Access must now be granted to ThetaHarvest.app (it is the responsible
       process for reaching the repo under ~/Downloads).
       LIVE mode (2026-06-05): no --shadow → writes the real history/ files (working tree only;
       the user reviews + pushes). To revert to validation, re-add <string>--shadow</string> below. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/victor_he/Downloads/Code/option-harvest/automation/launchd/ThetaHarvest.app/Contents/MacOS/Theta Harvest</string>
    <string>--quiet</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/victor_he/Downloads/Code/option-harvest</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- node (for claude) on PATH; CLAUDE_BIN absolute; ANTHROPIC_API_KEY deliberately absent -> Max plan -->
    <key>PATH</key>
    <string>/Users/victor_he/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>CLAUDE_BIN</key>
    <string>/Users/victor_he/.nvm/versions/node/v22.22.0/bin/claude</string>
  </dict>

  <!-- 09:00 Mac-local (SGT) daily; runs once on next wake if the Mac was asleep at 09:00. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/victor_he/Downloads/Code/option-harvest/automation/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/victor_he/Downloads/Code/option-harvest/automation/logs/launchd.err</string>
</dict>
</plist>
